# debian/lotos-screensaver.service
[Unit]
Description=Lotos Screensaver
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/local/bin/lotos-screensaver
Restart=always
