# debian/lotos-idle.path
[Unit]
Description=Watch lotos config.json

[Path]
PathChanged=%h/Documents/config/config.json

[Install]
WantedBy=default.target
